# startup-scripts/com.adrianstanca.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for Ollama (headless, no menu-bar app)
  Use this ONLY if Ollama.app is NOT installed.
  If Ollama.app is installed, it manages its own startup — skip this file.

  Install: cp com.adrianstanca.ollama.plist ~/Library/LaunchAgents/
           launchctl load ~/Library/LaunchAgents/com.adrianstanca.ollama.plist
  Check:   launchctl list | grep ollama
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.adrianstanca.ollama</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ollama</string>
    <string>serve</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Bind only to localhost for security -->
    <key>OLLAMA_HOST</key>
    <string>127.0.0.1:11434</string>
    <!-- Allow only 1 model loaded at a time to save VRAM -->
    <key>OLLAMA_MAX_LOADED_MODELS</key>
    <string>1</string>
    <!-- Single parallel request to avoid OOM with 27b model -->
    <key>OLLAMA_NUM_PARALLEL</key>
    <string>1</string>
    <!-- Keep model in memory for 10 min between runs -->
    <key>OLLAMA_KEEP_ALIVE</key>
    <string>10m</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/adrianstanca/.paperclip/logs/ollama.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/adrianstanca/.paperclip/logs/ollama.err.log</string>
</dict>
</plist>
