# startup-scripts/com.adrianstanca.paperclip.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for Paperclip
  Install: cp com.adrianstanca.paperclip.plist ~/Library/LaunchAgents/
           launchctl load ~/Library/LaunchAgents/com.adrianstanca.paperclip.plist
  Unload:  launchctl unload ~/Library/LaunchAgents/com.adrianstanca.paperclip.plist

  NOTE: Update ProgramArguments if your Paperclip binary path differs.
  Common locations:
    - /usr/local/bin/paperclipai   (global npm install)
    - /opt/homebrew/bin/paperclipai (Homebrew)
    - Run `which paperclipai` to confirm
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.adrianstanca.paperclip</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update this path: run `which paperclipai` in Terminal -->
    <string>/usr/local/bin/paperclipai</string>
    <string>dev</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/adrianstanca</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/adrianstanca</string>
    <key>OLLAMA_HOST</key>
    <string>http://localhost:11434</string>
  </dict>

  <!-- Start automatically when user logs in -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Throttle restarts: wait 10s before re-launching -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/adrianstanca/.paperclip/logs/paperclip.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/adrianstanca/.paperclip/logs/paperclip.err.log</string>
</dict>
</plist>
